# combust.service
[Unit]
Description=Combust Firewall Script

[Service]
Type=oneshot
ExecStart=/etc/iptables/combust.sh
ExecStop=/etc/iptables/combust.sh --flush
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
