# combust.service
[Unit]
Description=Netfilter (iptables) Firewall Script

[Service]
Type=oneshot
ExecStart=/etc/combust/combust.sh
ExecStop=/etc/combust/combust.sh --flush
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
